# REST/python/notebooks/Solaris_Zones/template-sc_profile.xml
<?xml version='1.0' encoding='UTF-8'?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!-- Auto-generated by sysconfig -->
<service_bundle type="profile" name="sysconfig">
  <service name="system/identity" type="service" version="1">
    <instance name="cert" enabled="true"/>
    <instance name="node" enabled="true">
      <property_group name="config" type="application">
        <propval name="nodename" value="{$nodename}" type="astring"/>
      </property_group>
    </instance>
  </service>
  <service name="network/install" type="service" version="1">
    <instance name="default" enabled="true">
      <property_group name="install_ipv4_interface" type="application">
        <propval name="name" value="net0/v4" type="astring"/>
        <propval name="address_type" value="static" type="astring"/>
        <propval name="static_address" value="{$IP_address}/24" type="net_address_v4"/>
        <propval name="default_route" value="{$IP_router}" type="net_address_v4"/>
      </property_group>
      <property_group name="install_ipv6_interface" type="application">
        <propval name="name" value="net0/v6" type="astring"/>
        <propval name="address_type" value="addrconf" type="astring"/>
        <propval name="stateful" value="yes" type="astring"/>
        <propval name="stateless" value="yes" type="astring"/>
      </property_group>
    </instance>
  </service>
  <service name="network/dns/client" type="service" version="1">
    <property_group name="config" type="application">
      <property name="nameserver" type="net_address">
        <net_address_list>
          <value_node value="8.8.8.8"/>
          <value_node value="8.8.4.4"/>
        </net_address_list>
      </property>
      <property name="search" type="astring">
        <astring_list>
          <value_node value="{$FQDN}"/>
        </astring_list>
      </property>
    </property_group>
    <instance name="default" enabled="true"/>
  </service>
  <service name="system/name-service/cache" type="service" version="1">
    <instance name="default" enabled="true"/>
  </service>
  <service name="system/name-service/switch" type="service" version="1">
    <property_group name="config" type="application">
      <propval name="default" value="files" type="astring"/>
      <propval name="host" value="files dns" type="astring"/>
    </property_group>
    <instance name="default" enabled="true"/>
  </service>
  <service name="system/environment" type="service" version="1">
    <instance name="init" enabled="true">
      <property_group name="environment" type="application">
        <propval name="LANG" value="C" type="astring"/>
      </property_group>
    </instance>
  </service>
  <service name="system/timezone" type="service" version="1">
    <instance name="default" enabled="true">
      <property_group name="timezone" type="application">
        <propval name="localtime" value="US/Pacific" type="astring"/>
      </property_group>
    </instance>
  </service>
  <service name="system/config-user" type="service" version="1">
    <instance name="default" enabled="true">
      <property_group name="user_account" type="application">
        <propval name="login" value="demo" type="astring"/>
        <propval name="password" value="{$user_passwd_hash}" type="astring"/>
        <propval name="description" value="demo" type="astring"/>
        <propval name="shell" value="/usr/bin/bash" type="astring"/>
        <propval name="gid" value="10" type="count"/>
        <propval name="auth_profiles" value="System Administrator" type="astring"/>
        <propval name="roles" value="root" type="astring"/>
        <propval name="sudoers" value="ALL=(ALL) ALL" type="astring"/>
        <propval name="type" value="normal" type="astring"/>
      </property_group>
      <property_group name="root_account" type="application">
        <propval name="login" value="root" type="astring"/>
        <propval name="password" value="{$root_passwd_hash}" type="astring"/>
        <propval name="type" value="role" type="astring"/>
      </property_group>
    </instance>
  </service>
</service_bundle>
